# configs/etc/local.d/lcd4linux.start
#!/bin/bash

modprobe cdc-acm
lcd4linux -q
