# configs/etc/local.d/noop_ssd.start
#!/bin/bash

echo noop > /sys/block/sda/queue/scheduler
